# web/cron.d/mafiascum
0 5 * * * /opt/bitnami/scripts/mafiascum/ms_perform_backup.sh