# database/cron.d/mafiascum
0 5 * * * root /opt/mafiascum/scripts/ms_perform_backup.sh
