# database/cron.d/mafiascum
0 5 * * * /opt/mafiascum/scripts/ms_perform_backup.sh